# conf/etc/periodic/hourly/01-mariadb
#!/usr/bin/with-contenv sh

DISABLE=$(cat /tmp/backuppart.disable)
if [ "$DISABLE" = "1" ];
then
    exit 0
fi
sleep  $(( $RANDOM % 30 +1 ))m
backup partial