# conf/etc/periodic/daily/01-mariadb
#!/usr/bin/with-contenv sh
dumps(){
    DB=$1
    PREFIX=$2
    echo -n "> "$DB
    ERROR=0
    mysqldump -uroot -p$ROOT_PASSWORD $DB > /var/lib/mysql-backup/$PREFIX$DB.sql
    if [ "$?" -eq 0 ]
    then
        echo " ........ OK"
    else
        echo " !!!!!!!! NOK"
        ERROR=1
    fi
    echo $ERROR  > /var/lib/mysql-backup/$PREFIX$DB.status
    
}

process(){
    dumps $1 "tmp-"
    if [ "$(cat /var/lib/mysql-backup/tmp-$1.status)" = "0" ];
    then
        mv /var/lib/mysql-backup/$1.sql.6 /var/lib/mysql-backup/$1.sql.7 > /dev/null 2>&1
        mv /var/lib/mysql-backup/$1.sql.5 /var/lib/mysql-backup/$1.sql.6 > /dev/null 2>&1
        mv /var/lib/mysql-backup/$1.sql.4 /var/lib/mysql-backup/$1.sql.5 > /dev/null 2>&1
        mv /var/lib/mysql-backup/$1.sql.3 /var/lib/mysql-backup/$1.sql.4 > /dev/null 2>&1
        mv /var/lib/mysql-backup/$1.sql.2 /var/lib/mysql-backup/$1.sql.3 > /dev/null 2>&1
        mv /var/lib/mysql-backup/$1.sql.1 /var/lib/mysql-backup/$1.sql.2 > /dev/null 2>&1
        mv /var/lib/mysql-backup/$1.sql /var/lib/mysql-backup/$1.sql.1 > /dev/null 2>&1
        mv /var/lib/mysql-backup/tmp-$1.sql /var/lib/mysql-backup/$1.sql > /dev/null 2>&1
        mv /var/lib/mysql-backup/tmp-$1.sql.status /var/lib/mysql-backup/$1.sql.status > /dev/null 2>&1
    fi
}


next(){
    echo "> Get Next Site"
    COUNTER=$(($COUNTER + 1))
}

is_db_enabled(){
    NAME="DB_"$COUNTER"_NAME"
    DB_NAME=$(eval 'echo $'$NAME)
    if [ -z "$DB_NAME" ];
    then
        return 1
    fi
    
    return 0
}


COUNTER=0
while [ 1 -eq 1 ];
do
    if ! is_db_enabled;
    then
        break
    fi
    process $DB_NAME
    next
done